# deploy/hans-webhook.service
[Unit]
Description=Hans GitHub Webhook Listener
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/hans
ExecStart=/home/pi/hans/.venv/bin/python deploy/update-hook.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
